# packer/SCRIPTS/lvm_ubuntu.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/lvm_ubuntu.sh

[Install]
WantedBy=default.target
